# schedule.cron
* * * * * /bin/sh /exporter/export.sh >> /exporter/logs/exporter.log 2>&1
